# Scripts/Editing_scripts/scripts/miniforge3/etc/profile.d/conda.csh
setenv CONDA_EXE "`cygpath 'C:/miniforge3\Scripts\conda.exe'`";
setenv _CONDA_ROOT "`cygpath 'C:/miniforge3'`";
setenv _CONDA_EXE "`cygpath 'C:/miniforge3\Scripts\conda.exe'`";
setenv CONDA_PYTHON_EXE "`cygpath 'C:/miniforge3\python.exe'`";
source "`cygpath 'C:/miniforge3\Lib\site-packages\conda\shell\etc\profile.d\conda.csh'`";
